# image_content/scripts/etc/profile.d/odp_session_variables.sh
export sessionusername=$USER
export sessiontoken=$(head -n 1 /var/lib/eric-odp/.enm_login)